# nasdaq-agent.service
[Unit]
Description=NASDAQ Stock Agent - AI-powered stock analysis service
Documentation=https://github.com/your-repo/nasdaq-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nasdaq-agent
Group=nasdaq-agent
WorkingDirectory=/opt/nasdaq-agent

# Load environment variables from .env file
EnvironmentFile=/opt/nasdaq-agent/.env

# Start the application
ExecStart=/opt/nasdaq-agent/venv/bin/python main.py

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nasdaq-agent

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
